# creatus.service
[Unit]
Description=
After=network.target

[Service]
Type=simple
ExecStart=/bin/creatus

[Install]
WantedBy=multi-user.target
